# usr/lib/udev/rules.d/99-eta-usb.rules
SUBSYSTEMS=="usb", RUN+="/usr/libexec/eta-usb-trigger"

