# etc/xdg/autostart/etausb-agent.desktop
[Desktop Entry]
Type=Application
Icon=gtk-ok
Name=Eta usb agent
Exec=/usr/libexec/etausb-agent
Comment=Eta Usb Login Client
