# etc/xdg/autostart/nemo-bookmarks.desktop
[Desktop Entry]
Name=Nemo Bookmarks
Exec=/usr/share/pardus/pardus-etap-settings/nemo-bookmarks
Type=Application
StartupNotify=false
NoDisplay=true
